# supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:zookeeper]
command=/usr/local/zookeeper/bin/zkServer.sh start-foreground
user=daemon
priority=0

[program:mysql]
command=/usr/bin/pidproxy /var/run/mysqld/mysqld.pid /usr/bin/mysqld_safe
user=mysql
priority=0

[program:druid-coordinator]
user=druid
command=java
  -server
  -Xmx1g
  -Duser.timezone=UTC
  -Dfile.encoding=UTF-8
  -Ddruid.host=%(ENV_HOSTIP)s
  -Ddruid.extensions.loadList=[\"mysql-metadata-storage\"]
  -Ddruid.extensions.directory=/usr/local/druid/extensions
  -Ddruid.extensions.hadoopDependenciesDir=/usr/local/druid/hadoop/hadoop_druid_dependencies
  -Ddruid.metadata.storage.type=mysql
  -Ddruid.metadata.storage.connector.connectURI=jdbc:mysql://localhost:3306/druid
  -Ddruid.metadata.storage.connector.user=druid
  -Ddruid.metadata.storage.connector.password=diurd
  -Ddruid.coordinator.startDelay=PT5S
  -Ddruid.worker.capacity=2
  -cp /usr/local/druid/lib/*
  io.druid.cli.Main server coordinator
redirect_stderr=true
priority=100

[program:druid-indexing-service]
user=druid
command=java
  -server
  -Xmx64m
  -Duser.timezone=UTC
  -Dfile.encoding=UTF-8
  -Ddruid.host=%(ENV_HOSTIP)s
  -Ddruid.extensions.loadList=[\"mysql-metadata-storage\"]
  -Ddruid.extensions.directory=/usr/local/druid/extensions
  -Ddruid.extensions.hadoopDependenciesDir=/usr/local/druid/hadoop/hadoop_druid_dependencies
  -Ddruid.metadata.storage.type=mysql
  -Ddruid.metadata.storage.connector.connectURI=jdbc:mysql://localhost:3306/druid
  -Ddruid.metadata.storage.connector.user=druid
  -Ddruid.metadata.storage.connector.password=diurd
  -Ddruid.indexer.storage.type=metadata
  -Ddruid.indexer.queue.startDelay=PT0M
  -Ddruid.indexer.runner.pendingTasksRunnerNumThreads=2
  -Ddruid.worker.capacity=4
  -Ddruid.indexer.runner.javaOpts="-server -Xmx2g -XX:+UseG1GC -XX:MaxDirectMemorySize=5294967296"
  -cp /usr/local/druid/lib/*
  io.druid.cli.Main server overlord

[program:druid-historical]
user=druid
command=java
  -server
  -Xmx1g
  -XX:MaxDirectMemorySize=3g
  -XX:+UseConcMarkSweepGC
  -XX:+PrintGCDetails
  -XX:+PrintGCTimeStamps
  -Duser.timezone=UTC
  -Dfile.encoding=UTF-8
  -Ddruid.historical.cache.useCache=true
  -Ddruid.historical.cache.populateCache=true
  -Ddruid.processing.buffer.sizeBytes=1073741824
  -Ddruid.processing.numThreads=7
  -Ddruid.server.http.numThreads=10
  -Ddruid.server.maxSize=300000000000
  -Ddruid.host=%(ENV_HOSTIP)s
  -Ddruid.extensions.loadList=[\"druid-s3-extensions\"]
  -Ddruid.extensions.directory=/usr/local/druid/extensions
  -Ddruid.extensions.hadoopDependenciesDir=/usr/local/druid/hadoop/hadoop_druid_dependencies
  -Ddruid.s3.accessKey=AKIAIMKECRUYKDQGR6YQ
  -Ddruid.s3.secretKey=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
  -Ddruid.computation.buffer.size=67108864
  -Ddruid.segmentCache.locations="[{\"path\":\"/var/tmp/druid/indexCache\",\"maxSize\":5000000000}]"
  -cp /usr/local/druid/lib/*
  io.druid.cli.Main server historical
redirect_stderr=true
priority=100

[program:druid-broker]
user=druid
command=java
  -server
  -Xmx1g
  -Duser.timezone=UTC
  -Dfile.encoding=UTF-8
  -Ddruid.processing.buffer.sizeBytes=1073741824
  -Ddruid.host=%(ENV_HOSTIP)s
  -Ddruid.computation.buffer.size=67108864
  -Ddruid.broker.cache.sizeInBytes=33554432
  -Ddruid.worker.capacity=4
  -Ddruid.broker.http.numConnections=5
  -Ddruid.processing.numThreads=7
  -Ddruid.server.http.numThreads=20
  -Ddruid.broker.http.readTimeout=PT5M
  -cp /usr/local/druid/lib/*
  io.druid.cli.Main server broker
redirect_stderr=true
priority=100

[program:tranquility]
user=druid
command=/usr/local/tranquility-distribution-0.8.0/bin/tranquility
  kafka -configFile /usr/local/tranquility-distribution-0.8.0/conf/kafka.json
redirect_stderr=true
priority=1000

[program:imply-pivot]
user=druid
command=pivot --druid localhost:8082
redirect_stderr=true
priority=10000

